# systemd/pakistan-scraper.service
[Unit]
Description=Pakistan Law Scraper Service
After=network.target network-online.target
Wants=network-online.target
Documentation=https://github.com/yourrepo/pakistan-law-scraper

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/pakistan_law_scraper
Environment="PATH=/opt/pakistan_law_scraper/venv/bin"
ExecStart=/opt/pakistan_law_scraper/venv/bin/python run_production.py start
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/pakistan-scraper.log
StandardError=append:/var/log/pakistan-scraper.log

# Graceful shutdown handling
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
